# infra/prod/system.d/cloner.timer
[Unit]
Description=rclone timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
